# crontab
* * * * * echo "Crontab output job 1" > /proc/1/fd/1
* * * * * echo "Crontab output job 2" > /proc/1/fd/1

