# crontab
* * * * * echo "Crontab is working - watchdog 1" > /proc/1/fd/1
* * * * * echo "Crontab is working second - watchdog 2" > /proc/1/fd/1

